# lobe/config/templates/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:{{ lobe_name }}-logica-web]
command={{ lobe_dir }}/env/bin/gunicorn -b 127.0.0.1:{{ webserver_port }} -w {{ gunicorn_workers }} --max-requests {{ gunicorn_max_requests }} --max-requests-jitter {{ gunicorn_max_requests_jitter }} -t {{ http_timeout }} logica.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/web.log
stderr_logfile={{ lobe_dir }}/logs/web.error.log
user={{ user }}
directory={{ sites_dir }}

{% if use_rq %}
[program:{{ lobe_name }}-logica-schedule]
command={{ lobe_cmd }} schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/schedule.log
stderr_logfile={{ lobe_dir }}/logs/schedule.error.log
user={{ user }}
directory={{ lobe_dir }}

[program:{{ lobe_name }}-logica-default-worker]
command={{ lobe_cmd }} worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/worker.log
stderr_logfile={{ lobe_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ lobe_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ lobe_name }}-logica-short-worker]
command={{ lobe_cmd }} worker --queue short
priority=4
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/worker.log
stderr_logfile={{ lobe_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=360
directory={{ lobe_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

[program:{{ lobe_name }}-logica-long-worker]
command={{ lobe_cmd }} worker --queue long
priority=4
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/worker.log
stderr_logfile={{ lobe_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1560
directory={{ lobe_dir }}
killasgroup=true
numprocs={{ background_workers }}
process_name=%(program_name)s-%(process_num)d

{% for worker_name, worker_details in workers.items() %}
[program:{{ lobe_name }}-logica-{{ worker_name }}-worker]
command={{ lobe_cmd }} worker --queue {{ worker_name }}
priority=4
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/worker.log
stderr_logfile={{ lobe_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs={{ worker_details["timeout"] }}
directory={{ lobe_dir }}
killasgroup=true
numprocs={{ worker_details["background_workers"] or background_workers }}
process_name=%(program_name)s-%(process_num)d
{% endfor %}

{% else %}
[program:{{ lobe_name }}-logica-workerbeat]
command={{ lobe_dir }}/env/bin/python -m logica.celery_app beat -s beat.schedule
priority=3
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/workerbeat.log
stderr_logfile={{ lobe_dir }}/logs/workerbeat.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ lobe_name }}-logica-worker]
command={{ lobe_dir }}/env/bin/python -m logica.celery_app worker -n jobs@%%h -Ofair --soft-time-limit 360 --time-limit 390 --loglevel INFO
priority=4
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/worker.log
stderr_logfile={{ lobe_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=400
directory={{ sites_dir }}
killasgroup=true

[program:{{ lobe_name }}-logica-longjob-worker]
command={{ lobe_dir }}/env/bin/python -m logica.celery_app worker -n longjobs@%%h -Ofair --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/worker.log
stderr_logfile={{ lobe_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1540
directory={{ sites_dir }}
killasgroup=true

[program:{{ lobe_name }}-logica-async-worker]
command={{ lobe_dir }}/env/bin/python -m logica.celery_app worker -n async@%%h -Ofair --soft-time-limit 1500 --time-limit 1530 --loglevel INFO
priority=2
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/worker.log
stderr_logfile={{ lobe_dir }}/logs/worker.error.log
user={{ user }}
stopwaitsecs=1540
directory={{ sites_dir }}
killasgroup=true

{% endif %}

{% if not skip_redis %}
[program:{{ lobe_name }}-redis-cache]
command={{ redis_server }} {{ redis_cache_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/redis-cache.log
stderr_logfile={{ lobe_dir }}/logs/redis-cache.error.log
user={{ user }}
directory={{ sites_dir }}

[program:{{ lobe_name }}-redis-queue]
command={{ redis_server }} {{ redis_queue_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/redis-queue.log
stderr_logfile={{ lobe_dir }}/logs/redis-queue.error.log
user={{ user }}
directory={{ sites_dir }}
{% endif %}

{% if not skip_redis %}
[program:{{ lobe_name }}-redis-socketio]
command={{ redis_server }} {{ redis_socketio_config }}
priority=1
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/redis-socketio.log
stderr_logfile={{ lobe_dir }}/logs/redis-socketio.error.log
user={{ user }}
directory={{ sites_dir }}
{% endif %}

{% if node %}
[program:{{ lobe_name }}-node-socketio]
command={{ node }} {{ lobe_dir }}/apps/logica/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile={{ lobe_dir }}/logs/node-socketio.log
stderr_logfile={{ lobe_dir }}/logs/node-socketio.error.log
user={{ user }}
directory={{ lobe_dir }}
{% endif %}

[group:{{ lobe_name }}-web]
programs={{ lobe_name }}-logica-web {%- if node -%} ,{{ lobe_name }}-node-socketio {%- endif%}

{% if use_rq %}

[group:{{ lobe_name }}-workers]
programs={{ lobe_name }}-logica-schedule,{{ lobe_name }}-logica-default-worker,{{ lobe_name }}-logica-short-worker,{{ lobe_name }}-logica-long-worker

{% else %}

[group:{{ lobe_name }}-workers]
programs={{ lobe_name }}-logica-workerbeat,{{ lobe_name }}-logica-worker,{{ lobe_name }}-logica-longjob-worker,{{ lobe_name }}-logica-async-worker

{% endif %}

{% if not skip_redis %}
[group:{{ lobe_name }}-redis]
programs={{ lobe_name }}-redis-cache,{{ lobe_name }}-redis-queue,{{ lobe_name }}-redis-socketio
{% endif %}
