# lobe/config/templates/systemd/frappe-bench-web.target
[Unit]
After=network.target
Wants={{ lobe_name }}-logica-web.service {{ lobe_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
